# roles/result/files/result.service
[Unit]
Description=Result Web Client
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/result
Environment="DB_HOST=postgresql-1"
Environment="DB_NAME=paul"
Environment="DB_USER=paul"
Environment="DB_PASSWORD=democracyIsFragile"
ExecStart=/opt/result/venv/bin/python app.py
Restart=always

[Install]
WantedBy=multi-user.target
je trouvre les